# tools/saltstead-dash.service
[Unit]
Description=Saltstead harbourmaster ledger (invite codes + login tokens)
After=network-online.target

[Service]
Type=simple
User=james
WorkingDirectory=/home/james/saltstead/dash
Environment=PYTHONUTF8=1
ExecStart=/home/james/moorstead/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8097 --log-level warning
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
